# etc/systemd/system/bluebinder.service.d/override.conf
# Override service bluebinder
[Service]
ExecStartPre=/bin/sleep 5
